# deploy/ppocr-server.service
[Unit]
Description=PP-OCRv6 OCR server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Create a dedicated service user, or replace these with the account that owns
# /opt/ocr-servos and its model cache.
User=ppocr
Group=ppocr
WorkingDirectory=/opt/ocr-servos
Environment=RUST_LOG=info
# Models are fetched from the CDN on first run and cached here (persist across
# restarts so they download only once). Needs network on the very first start.
Environment=PPOCR_CACHE_DIR=/opt/ocr-servos/models-cache
ExecStart=/opt/ocr-servos/ppocr-server --host 127.0.0.1 --port 3088 --plan-cache 16
Restart=on-failure
RestartSec=3
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
